# cron.daily/update
#!/bin/sh
cd /home/alan/.dotfiles && make update
exit 0
